# roles/mail_server/files/spam-training.timer
[Unit]
Description=Run spam-training.service hourly.

[Timer]
OnCalendar=hourly
Persistent=true

