# roles/postgresql_server/files/postgresql-dumpall.timer
[Unit]
Description=Run postgresql-dumpall.service hourly.

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
