# roles/tarsnap/files/tarsnap-backups.timer
[Unit]
Description=Run tarsnap-backups.service daily.

[Timer]
OnCalendar=daily
Persistent=true

