# roles/zfs_pool/templates/zfs-ssd-pool-scrub.service.j2
[Unit]
Description=ZFS Scrub on ssd_pool

[Service]
Type=oneshot
ExecStart=/sbin/zpool scrub ssd_pool
