# roles/zfs_pool/templates/zfs-ssd-pool-healthcheck.timer.j2
[Unit]
Description=Weekly ZFS Health Check Timer for ssd_pool

[Timer]
OnCalendar={{ healthcheck_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
