# roles/spideroak/files/postgresql-dumpall.timer
[Unit]
Description=Run postgresql-dumpall.service hourly.

[Timer]
OnCalendar=hourly
Persistent=true
