# roles/spideroak/files/spideroak-backups.timer
[Unit]
Description=Run spideroak-backups.service daily.

[Timer]
OnCalendar=daily
Persistent=true

